# nmap.crontab
# m h  dom mon dow   command
0 3 * * * /root/nmap/nightly-nmap